# systemd/alpha-live.timer
[Unit]
Description=Run Sharpe-Renaissance Alpha Live Cycle daily

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=1800

[Install]
WantedBy=timers.target

