# systemd/alpha-scorecard.timer
[Unit]
Description=Run Sharpe-Renaissance Alpha Paper Scorecard daily

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=2400

[Install]
WantedBy=timers.target

