# systemd/reddit-ingest.timer
[Unit]
Description=Run Sharpe-Renaissance Reddit ingestion daily

[Timer]
OnCalendar=*-*-* 18:10:00 UTC
Persistent=true
RandomizedDelaySec=10m
Unit=reddit-ingest.service

[Install]
WantedBy=timers.target

